# bytd/bytd.socket
[Socket]
ListenStream=1981

[Install]
WantedBy=sockets.target
